# deploy/systemd/open-transit-scorecard-export.timer
[Unit]
Description=Run Open Transit RT pilot scorecard export hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

