# deploy/systemd/open-transit-feed-monitor.service
[Unit]
Description=Open Transit RT pilot public feed monitor
Documentation=https://github.com/YOUR_ORG/open-transit-rt
After=network.target

[Service]
Type=oneshot
User={{OCI_APP_USER}}
Group={{OCI_APP_USER}}
WorkingDirectory={{OCI_REMOTE_DIR}}/app
EnvironmentFile={{OCI_REMOTE_DIR}}/ops/pilot-ops.env
ExecStart={{OCI_REMOTE_DIR}}/app/scripts/pilot-ops.sh feed-monitor
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{OCI_REMOTE_DIR}}

