# deploy/systemd/open-transit-telemetry-ingest.service
[Unit]
Description=Open Transit RT — telemetry-ingest
Documentation=https://github.com/YOUR_ORG/open-transit-rt
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User={{OCI_APP_USER}}
WorkingDirectory={{OCI_REMOTE_DIR}}/app
EnvironmentFile={{OCI_REMOTE_DIR}}/env
Environment="PORT=8082"
ExecStart={{OCI_REMOTE_DIR}}/bin/telemetry-ingest
Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{OCI_REMOTE_DIR}}

[Install]
WantedBy=multi-user.target
