# deploy/systemd/open-transit-validator-cycle.timer
[Unit]
Description=Run Open Transit RT pilot validator cycle hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

