# deploy/systemd/open-transit-backup.timer
[Unit]
Description=Run Open Transit RT pilot backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

